# setup/systemd/noxflow-weekly-health.service
[Unit]
Description=Noxflow weekly health check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/Documents/code/dotfiles/setup/weekly-health-check.sh
