# systemd/user/wayland-wm@hyprland.desktop.service.d/10-aq-drm-devices.conf
[Service]
# Keep AQ_DRM_DEVICES dynamic via ~/.config/uwsm/env-hyprland.
# A static card order can flip across boots and break Hyprland startup.
UnsetEnvironment=AQ_DRM_DEVICES
